# oshd-clean.service
# oshd-noclear
# Only install if using a package manager. Otherwise use oshctl service install.
[Unit]
Description=oshd Cleanup
DefaultDependencies=no
Before=shutdown.target reboot.target halt.target poweroff.target

[Service]
Type=oneshot
ExecStart=/usr/bin/oshctl clear

[Install]
WantedBy=halt.target reboot.target shutdown.target poweroff.target
